# supervisord.conf
[supervisord]
nodaemon=true

[program:httpd]
command=/bin/bash -c "rm -rf /run/apache2/* && /usr/sbin/apachectl -D FOREGROUND"
